# autoUpdateCode.ini
[program:autoUpdateCode]
process_name=autoUpdateCode
command=/usr/bin/sh /data/wwwroot/eliveapi/autoUpdateCode.sh eliveapi
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/data/wwwlogs/autoUpdateCode-supervisor.log
